# scripts/systemd/invest-pi-outcomes.timer
[Unit]
Description=Invest-Pi · Outcomes täglich 02:30

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
