# scripts/systemd/invest-pi-rebalance.service
[Unit]
Description=Invest-Pi · Rebalance/Sell-Pass (Mo-Fr 20:00 CET, sell-only)
After=invest-pi-score.service

[Service]
Type=oneshot
User=investpi
WorkingDirectory=/home/investpi/invest-pi
EnvironmentFile=/home/investpi/invest-pi/.env
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/usr/bin/python3 -B scripts/run_strategy.py --skip-buys
MemoryMax=2G
CPUQuota=50%
