# scripts/systemd/invest-pi-hardware.service
[Unit]
Description=Invest-Pi · Hardware-Check (CPU/Disk/Mem) mit Telegram-Push

OnFailure=invest-pi-error-alert@%n.service
[Service]
Type=oneshot
User=investpi
WorkingDirectory=/home/investpi/invest-pi
EnvironmentFile=/home/investpi/invest-pi/.env
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/usr/bin/python3 -B scripts/hardware_check.py
MemoryMax=128M
CPUQuota=10%
