# scripts/systemd/invest-pi-dca-watchdog.service
[Unit]
Description=Invest-Pi · DCA Watchdog (daily position monitoring)

[Service]
Type=oneshot
User=investpi
WorkingDirectory=/home/investpi/invest-pi
EnvironmentFile=/home/investpi/invest-pi/.env
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/usr/bin/python3 -B scripts/dca_watchdog.py
MemoryMax=256M
CPUQuota=20%
