# scripts/systemd/invest-pi-train-regime.timer
[Unit]
Description=Invest-Pi · HMM-Training Sa 05:00 CEST

[Timer]
OnCalendar=Sat *-*-* 05:00:00
Persistent=true

[Install]
WantedBy=timers.target
