# scripts/systemd/invest-pi-ai-swing-outcomes.timer
[Unit]
Description=Invest-Pi · KI-Swing Outcome-Tracker taeglich 03:05

[Timer]
OnCalendar=*-*-* 03:05:00
Persistent=true

[Install]
WantedBy=timers.target
