# scripts/systemd/invest-pi-daily-report.service
[Unit]
Description=Invest-Pi · Daily Performance Report (Telegram)

[Service]
Type=oneshot
User=investpi
WorkingDirectory=/home/investpi/invest-pi
EnvironmentFile=/home/investpi/invest-pi/.env
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/usr/bin/python3 -B scripts/daily_report.py --mode auto
MemoryMax=256M
CPUQuota=20%
