# scripts/systemd/invest-pi-meta-review.timer
[Unit]
Description=Invest-Pi · Meta-Review monatlich am 2. um 04:00 CEST

[Timer]
OnCalendar=*-*-02 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
