# scripts/systemd/invest-pi-sync.service
[Unit]
Description=Invest-Pi · Broker→DB Sync + Equity-Snapshot
After=network-online.target

[Service]
Type=oneshot
User=investpi
WorkingDirectory=/home/investpi/invest-pi
EnvironmentFile=/home/investpi/invest-pi/.env
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/usr/bin/python3 -B scripts/sync_positions.py
MemoryMax=1G
CPUQuota=30%
