# scripts/systemd/invest-pi-score.timer
[Unit]
Description=Invest-Pi · Risk-Scoring (stündlich :30)

[Timer]
OnCalendar=*-*-* *:30:00
Persistent=true

[Install]
WantedBy=timers.target
