# scripts/systemd/invest-pi-strategy.timer
[Unit]
Description=Invest-Pi · Strategy-Run Mo-Fr 18:00 (Midday)

[Timer]
OnCalendar=Mon..Fri *-*-* 18:00:00
Persistent=true

[Install]
WantedBy=timers.target
